# stuff/sinit/rc.local
#!/bin/sh
#
# User's commands started at boot
#
