# packages/lvm2/lvm2.install
post_install() {
	echo "Make sure the \"lvm\" init script is in the runlevels:"
	echo "# rc-update add lvm boot"
	echo
	echo "Make sure to enable lvmetad in /etc/lvm/lvm.conf if you want"
	echo "to enable lvm autoactivation and metadata caching."
}
