# packages/kmod/kmod.install
post_install() {
	rc-update add kmod-static-nodes sysinit
}

post_remove() {
	rc-update del kmod-static-nodes sysinit
}

post_upgrade() {
	post_install || true
}
