# KEEP/etc/inittab
# /etc/inittab

::sysinit:/etc/rc.boot
::respawn:/bin/runsvdir -P /var/service
::shutdown:/etc/rc.shutdown
::ctrlaltdel:/sbin/reboot
