# KEEP/etc/rc.local
#!/bin/sh

export PATH="/local/bin:/bin"
