# KEEP/etc/service/tty3/run
#!/bin/sh
exec getty -L 38400 tty3 linux
