# packages/eudev/eudev.install
post_install() {
	rc-update add udev sysinit
	rc-update add udev-trigger sysinit
	rc-update add udev-settle sysinit
	rc-update add udev-postmount sysinit
}

post_remove() {
	rc-update del udev sysinit
	rc-update del udev-trigger sysinit
	rc-update del udev-settle sysinit
	rc-update del udev-postmount sysinit
}

post_upgrade() {
	post_install || true
}
