# packages/eudev/eudev.install
post_install() {
	rc-update add udev sysinit
	rc-update add udev-trigger sysinit
}

post_remove() {
	rc-update del udev sysinit
	rc-update del udev-trigger sysinit
}

pre_upgrade() {
	post_remove
}

post_upgrade() {
	post_install
}
