# packages/base-files/inittab
# /etc/inittab

::sysinit:/bin/openrc sysinit
::sysinit:/bin/openrc boot
::wait:/bin/openrc default

tty1::respawn:/bin/getty 38400 tty1
tty2::respawn:/bin/getty 38400 tty2
tty3::respawn:/bin/getty 38400 tty3
tty4::respawn:/bin/getty 38400 tty4
tty5::respawn:/bin/getty 38400 tty5
tty6::respawn:/bin/getty 38400 tty6
#ttyS0::respawn:/bin/getty -L ttyS0 115200 vt100

::ctrlaltdel:/bin/reboot

::shutdown:/bin/openrc shutdown
