# KEEP/etc/service/ttyS0/run
#!/bin/sh
exec getty -L 38400 ttyS0 linux
