# stuff/base-files/inittab
::sysinit:/etc/rc.d/rc.boot
::respawn:/usr/bin/runsvdir-start
::shutdown:/etc/rc.d/rc.down
::ctrlaltdel:/usr/bin/reboot
