# KEEP/base-files/inittab
# /etc/inittab

::sysinit:/etc/rc
::respawn:/bin/runsvdir -P /var/service
::shutdown:/etc/rc.sysdown
::ctrlaltdel:/sbin/reboot
