# stuff/rc.d/rc.local
#!/bin/sh
#
# Put your commands on boot
#
