# stuff/svc/acpid.service
[Unit]
Description=Busybox ACPI daemon

[Service]
ExecStart=/usr/bin/acpid -g wheel -k 027

[Install]
WantedBy=multi-user.target
