# stuff/svc/crond.service
[Unit]
Description=Busybox cron daemon

[Service]
ExecStart=/usr/bin/crond -f -d8

[Install]
WantedBy=multi-user.target
