# stuff/svc/dhcrelay.service
[Unit]
Description=DHCP Relay Agent Daemon
Documentation=man:dhcrelay(8)
Wants=network-online.target
After=network-online.target

[Service]
Type=notify
ExecStart=/usr/bin/dhcrelay -d --no-pid
StandardError=null

[Install]
WantedBy=multi-user.target
