# stuff/svc/dnsmasq.service
[Unit]
Description=DNS caching server.
After=network.target

[Service]
ExecStart=/usr/bin/dnsmasq -k

[Install]
WantedBy=multi-user.target
