# stuff/svc/inetd.service
[Unit]
Description=Busybox INET daemon

[Service]
ExecStart=/usr/bin/inetd

[Install]
WantedBy=multi-user.target
