# stuff/svc/krb5kdc.service
[Unit]
Description=Kerberos 5 KDC
Wants=network-online.target
After=syslog.target network.target network-online.target

[Service]
Type=forking
PIDFile=/run/krb5kdc.pid
EnvironmentFile=-/etc/conf.d/krb5kdc
ExecStart=/usr/bin/krb5kdc -P /run/krb5kdc.pid $KRB5KDC_ARGS
ExecReload=/usr/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
