# stuff/svc/pppoe-server.service
[Unit]
Description=PPPoE Server.
After=syslog.target

[Service]
ExecStart=/usr/bin/pppoe-server

[Install]
WantedBy=multi-user.target
