# stuff/svc/rdate.service
[Unit]
Description=Busybox rdate daemon

[Service]
ExecStart=/usr/bin/rdate

[Install]
WantedBy=multi-user.target
