# stuff/svc/watchdog.service
[Unit]
Description=Busybox WatchDog daemon

[Service]
ExecStart=/usr/bin/watchdog /dev/watchdog -F

[Install]
WantedBy=multi-user.target
