# packages/filesystem/inittab
# /etc/inittab

::sysinit:/usr/bin/rc.init

::respawn:/usr/bin/runsvdir-start

::shutdown:/usr/bin/rc.shutdown
::ctrlaltdel:/sbin/reboot
