# packages/filesystem/rc.local
#!/bin/sh

# Append our default paths
appendpath () {
    case ":$PATH:" in
        *:"$1":*)
            ;;
        *)
            PATH="${PATH:+$PATH:}$1"
    esac
}

appendpath '/usr/bin'
appendpath '/usr/local/bin'
appendpath '/usr/local/sbin'
unset appendpath

export PATH
