# stuff/rc.local
#!/bin/sh

#dmesg >/var/log/dmesg
#syslogd -D -f /etc/syslog.conf
#ntpd -s -d
#sshd -De
