# stuff/rc.local
#!/bin/sh

PATH="/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/janus/bin:/janus/sbin"
LD_LIBRARY_PATH="/usr/lib:/lib:/janus/lib"

export PATH LD_LIBRARY_PATH

dmesg >/var/log/dmesg
syslogd -D -f /etc/syslog.conf
ntpd -s -d
sshd -De
